# systemd/pixelboard-watcher.path
[Path]
PathChanged=/app/pixelboard/pixelboard

[Install]
WantedBy=multi-user.target
